# deploy/systemd/geist-home.service
# Example only. The installer generates a host-specific unit from --binary,
# --work-dir, --user, and --ha-config; prefer that generated unit.
#
# Install after adapting the paths and user:
#   sudo cp geist-home.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now geist-home
#
# The socket lives in the Home Assistant config volume so the HA container
# (and its geist_conversation component) reaches it at /config/geist.sock.

[Unit]
Description=geist home appliance daemon (--serve)
After=local-fs.target

[Service]
User=geist
WorkingDirectory=/var/lib/geist
Environment=GEIST_AGENT_TRACE=0
ExecStart=/usr/local/bin/geist --serve /srv/homeassistant/geist.sock
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
